# roles/swap/files/ephemeral-swap.service
[Unit]
Description=Mount ephemeral swap partition

[Service]
Type=oneshot
ExecStart=/opt/swap.sh

[Install]
WantedBy=multi-user.target
